# rad.conf
[program:rad]
command=/home/rfo/wa/rad/ve/bin/python rad.py
directory=/home/rfo/wa/rad
autostart=true
autorestart=true
stderr_logfile=/var/log/rad.log
stdout_logfile=/var/log/rad.err
user = rfo
